# root/etc/services.d/vhosts-watch/run
#!/usr/bin/with-contenv sh

while inotifywait -e modify /root/.config/clash; do
    s6-svc -r /var/run/s6/services/clash
done